# puppetmaster/etc/s6/unicorn/run
#!/bin/sh

exec /usr/lib/ruby/gems/2.1.0/bin/unicorn -c /etc/puppet/unicorn.conf
