# puppetboard/files/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
command=/usr/sbin/uwsgi --master --ini /etc/nginx/puppetboard.ini
redirect_stderr=true
